# execution/openclaw-job.service
[Unit]
Description=OpenClaw Job Search Automation
After=network.target

[Service]
Type=oneshot
User=azureuser
WorkingDirectory=/home/azureuser/openclaw-job-hunter
ExecStart=/usr/bin/node execution/job-search.js
StandardOutput=append:/home/azureuser/openclaw-job-hunter/logs/job-search.log
StandardError=append:/home/azureuser/openclaw-job-hunter/logs/job-search-error.log

[Install]
WantedBy=multi-user.target
